# pipeline/desktop/deploy/desktop-pipeline.timer
[Unit]
Description=Run the Tomboy diary OCR pipeline every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Unit=desktop-pipeline.service

[Install]
WantedBy=timers.target
