# bridge/deploy/term-bridge.container
# Podman Quadlet unit for the tomboy-web terminal bridge.
#
# Install:
#   1. Build the image:
#        cd bridge && podman build -t term-bridge:latest .
#   2. Drop this file at ~/.config/containers/systemd/term-bridge.container
#   3. Create ~/.config/term-bridge.env with the required env vars:
#        BRIDGE_PASSWORD=...
#        BRIDGE_SECRET=...        (32+ random bytes, e.g. `openssl rand -hex 32`)
#        BRIDGE_ALLOWED_ORIGIN=https://your-app-domain
#        # Ollama 가 같은 호스트면 생략 가능 (코드 default 가 localhost:11434).
#        # 분리된 머신이면 (예: bridge=Pi, Ollama=Desktop) 명시 필요:
#        #   OLLAMA_BASE_URL=http://<ollama-host-ip>:11434
#        # RAG 검색 서버 (desktop FastAPI). 같은 호스트면 생략 가능.
#        # 분리된 머신이면 (예: bridge=Pi, indexer=Desktop) 명시:
#        #   RAG_SEARCH_URL=http://<desktop-host-ip>:8743/search
#        # OCR 서비스 URL. 같은 머신이면 http://localhost:8080. 분리된 머신
#        # (예: bridge=Pi, ocr=Desktop) 면 명시 필요:
#        #   OCR_SERVICE_URL=http://<desktop-host-ip>:8080
#        # 누락 시 부팅 거부. 같은 머신 가정으로 빠뜨려서 OCR이 안 되는
#        # 실수를 사전 차단한다.
#        BRIDGE_PUBLIC_BASE_URL=https://your-bridge-host
#        # The Caddy-fronted HTTPS hostname operators reach the bridge at.
#        # Embedded in /files upload responses; clients use it as the
#        # download URL.
#        # 누락 시 부팅 거부.
#   4. Create ~/.config/term-bridge/hosts.json BEFORE starting the unit
#      (podman creates a directory at the source path if the file is
#      missing, which silently breaks the bind mount). For the no-WOL
#      case write `{}`; for WOL hosts:
#        {
#          "desktop.lan": {
#            "mac": "AA:BB:CC:DD:EE:FF",
#            "broadcast": "192.168.0.255",
#            "wakeTimeoutSec": 90
#          }
#        }
#      Keys are matched against the host token in the note's ssh:// URL
#      (case-insensitive). Hosts not in the file = WOL skipped, ssh
#      attempted directly (existing behaviour).
#   4a. Pre-create the file storage directory BEFORE starting the unit:
#        mkdir -p ~/.local/share/term-bridge/files
#      As with hosts.json, podman creates a FILE at the source path if it
#      doesn't exist, which silently breaks the bind mount. The bridge
#      refuses to boot without BRIDGE_FILES_DIR pointing at a writable
#      directory.
#   5. (선택) 리마커블 배경화면 기능을 쓰면 ~/.config/term-bridge/remarkable.json
#      을 만든다(없거나 유효 항목이 없으면 /remarkable/wallpaper 는 503).
#      별칭 → 접속 좌표:
#        {
#          "rm2": { "host": "10.0.0.42", "user": "root", "port": 22 }
#        }
#      "host" 의 키(별칭)는 노트의 remarkable://<별칭> 과 매칭된다. SSH
#      인증은 위 ~/.ssh 마운트의 키를 그대로 쓰므로, 사전에 브릿지
#      호스트에서 리마커블로 `ssh-copy-id root@<리마커블-IP>` 로 키를 한 번
#      등록할 것.
#      또한 ~/.ssh 는 읽기전용 마운트라 첫 접속 시 host key 를
#      known_hosts 에 못 적는다 — `ssh-keyscan -H <리마커블-IP> >>
#      ~/.ssh/known_hosts` 로 미리 채워 둘 것.
#   5b. (선택) 폰/원격 역터널 별칭을 쓰면 ~/.config/term-bridge/ssh-hosts.json
#       을 만든다. 별칭 → 접속 좌표:
#         { "phone": { "host": "localhost", "port": 18022, "user": "termux" } }
#       노트의 ssh://<별칭> 과 매칭된다. 파일이 없으면 별칭 비활성(부팅 정상).
#       remarkable.json 과 마찬가지로 podman이 결측 source를 디렉토리로
#       만들어 마운트를 깨뜨리므로, 유닛 시작 前에 파일을 먼저 만들 것.
#   6. systemctl --user daemon-reload
#      systemctl --user enable --now term-bridge.service
#   7. Make sure linger is enabled so the container survives logout:
#        loginctl enable-linger $USER

[Unit]
Description=tomboy-web terminal bridge
After=network-online.target
Wants=network-online.target

[Container]
Image=localhost/term-bridge:latest
ContainerName=term-bridge
Network=host
UserNS=keep-id
EnvironmentFile=%h/.config/term-bridge.env
# Mount the user's SSH config + keys read-only so `ssh other-host` from
# inside the container uses the same agent/known_hosts the shell would.
Volume=%h/.ssh:/home/node/.ssh:ro,z
# WOL host map — read-only. File can be absent; bridge then logs "WOL
# disabled" and skips wake for every target.
Volume=%h/.config/term-bridge/hosts.json:/etc/term-bridge/hosts.json:ro,z
Environment=BRIDGE_HOSTS_FILE=/etc/term-bridge/hosts.json
# 워크로그 노트 creds ({uid, notebook, serviceAccount}) — 없으면 /notes는 503.
# 주의: podman은 결측 source를 디렉토리로 생성함 — 유닛 start 전에 파일을 먼저 만들 것.
Volume=%h/.config/term-bridge/notes-creds.json:/etc/term-bridge/notes-creds.json:ro,z
Environment=BRIDGE_NOTES_FILE=/etc/term-bridge/notes-creds.json
# 노트 changeDate는 Tomboy 로컬-offset 문자열 — 앱(폰, +09:00)과 offset이 어긋나면
# 문자열 LWW/정렬이 시간순과 뒤집힌다. 컨테이너 기본 UTC 방지용으로 앱과 동일 tz 고정.
Environment=TZ=Asia/Seoul
# reMarkable wallpaper host map — read-only, alias → {host,user,port?}.
# When the file is absent or has no valid host entries, the bridge logs
# that wallpaper is disabled and /remarkable/wallpaper returns 503. SSH
# auth reuses the read-only ~/.ssh mount above.
Volume=%h/.config/term-bridge/remarkable.json:/etc/term-bridge/remarkable.json:ro,z
Environment=BRIDGE_REMARKABLE_HOSTS_FILE=/etc/term-bridge/remarkable.json
# SSH 별칭 맵 — read-only, alias → {host,user?,port?}. 노트의 ssh://<별칭>
# 을 실제 접속 좌표로 치환한다(예: 역터널 phone → localhost:18022). 파일이
# 없거나 유효 항목이 없으면 별칭 비활성(부팅은 정상, 일반 ssh:// 만 동작).
Volume=%h/.config/term-bridge/ssh-hosts.json:/etc/term-bridge/ssh-hosts.json:ro,z
Environment=BRIDGE_SSH_HOSTS_FILE=/etc/term-bridge/ssh-hosts.json
# Persistent storage for /files uploads. Pre-create the directory on the
# host (see README install step) — podman creates a file at the source
# path otherwise, silently breaking the bind mount.
Volume=%h/.local/share/term-bridge/files:/var/lib/term-bridge/files:rw,z
Environment=BRIDGE_FILES_DIR=/var/lib/term-bridge/files
# 일기 파이프라인 드릴다운(GET /status/diary)용 — diary-sync inbox 읽기전용.
# 공유 디렉터리이므로 z/Z 리라벨 금지(diary-sync 접근 보존). 호스트 유저가 inbox 를
# 읽을 수 있어야 함(ACL traverse + inbox drwxr-xr-x — 2026-06-21 확인).
Volume=/home/diary-sync/diary/inbox:/var/lib/diary-inbox:ro
Environment=DIARY_INBOX_DIR=/var/lib/diary-inbox
# desktop trigger_server 가 도달 가능하면 OCR 섹션 활성(없으면 unconfigured 로 생략).
# Environment=DIARY_TRIGGER_URL=http://<desktop-host>:8765
# DIARY_TRIGGER_TOKEN 은 시크릿이므로 ~/.config/term-bridge.env 에 둔다.
# 리마커블 inbox 볼륨/SSH 환경변수는 더 이상 필요 없다.
# `리마커블::` 노트의 /remarkable/upload 는 단순히 automation-service의
# `pipeline-run` 을 깨우는 트리거이고, 리마커블→Pi inbox rsync 는 리마커블
# 위 `diary-push.timer`(1분 주기, mtime 가드)가 자율적으로 처리한다.
# Public URL the bridge returns to clients in upload responses. Must be
# the Caddy-fronted HTTPS hostname.
# Set in ~/.config/term-bridge.env: BRIDGE_PUBLIC_BASE_URL=https://term.example.com
# OLLAMA_BASE_URL 는 의도적으로 여기 hardcode 하지 않음 — EnvironmentFile
# 의 값이 multi-host setup (예: bridge=Pi, Ollama=Desktop) 에서 override
# 할 수 있도록. 코드 default 는 localhost:11434 (lib/llm.ts).
ReadOnly=true
NoNewPrivileges=true
DropCapability=ALL
Tmpfs=/tmp:rw,size=64m
# Reverse-proxied by Caddy/nginx on the host; no public bind here.

[Service]
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
