# bridge/deploy/term-bridge.container
# Podman Quadlet unit for the tomboy-web terminal bridge.
#
# Install:
#   1. Build the image:
#        cd bridge && podman build -t term-bridge:latest .
#   2. Drop this file at ~/.config/containers/systemd/term-bridge.container
#   3. Create ~/.config/term-bridge.env with the required env vars:
#        BRIDGE_PASSWORD=...
#        BRIDGE_SECRET=...        (32+ random bytes, e.g. `openssl rand -hex 32`)
#        BRIDGE_ALLOWED_ORIGIN=https://your-app-domain
#        # Ollama 가 같은 호스트면 생략 가능 (코드 default 가 localhost:11434).
#        # 분리된 머신이면 (예: bridge=Pi, Ollama=Desktop) 명시 필요:
#        #   OLLAMA_BASE_URL=http://<ollama-host-ip>:11434
#   4. Create ~/.config/term-bridge/hosts.json BEFORE starting the unit
#      (podman creates a directory at the source path if the file is
#      missing, which silently breaks the bind mount). For the no-WOL
#      case write `{}`; for WOL hosts:
#        {
#          "desktop.lan": {
#            "mac": "AA:BB:CC:DD:EE:FF",
#            "broadcast": "192.168.0.255",
#            "wakeTimeoutSec": 90
#          }
#        }
#      Keys are matched against the host token in the note's ssh:// URL
#      (case-insensitive). Hosts not in the file = WOL skipped, ssh
#      attempted directly (existing behaviour).
#   5. systemctl --user daemon-reload
#      systemctl --user enable --now term-bridge.service
#   6. Make sure linger is enabled so the container survives logout:
#        loginctl enable-linger $USER

[Unit]
Description=tomboy-web terminal bridge
After=network-online.target
Wants=network-online.target

[Container]
Image=localhost/term-bridge:latest
ContainerName=term-bridge
Network=host
UserNS=keep-id
EnvironmentFile=%h/.config/term-bridge.env
# Mount the user's SSH config + keys read-only so `ssh other-host` from
# inside the container uses the same agent/known_hosts the shell would.
Volume=%h/.ssh:/home/node/.ssh:ro,z
# WOL host map — read-only. File can be absent; bridge then logs "WOL
# disabled" and skips wake for every target.
Volume=%h/.config/term-bridge/hosts.json:/etc/term-bridge/hosts.json:ro,z
Environment=BRIDGE_HOSTS_FILE=/etc/term-bridge/hosts.json
# OLLAMA_BASE_URL 는 의도적으로 여기 hardcode 하지 않음 — EnvironmentFile
# 의 값이 multi-host setup (예: bridge=Pi, Ollama=Desktop) 에서 override
# 할 수 있도록. 코드 default 는 localhost:11434 (lib/llm.ts).
ReadOnly=true
NoNewPrivileges=true
DropCapability=ALL
Tmpfs=/tmp:rw,size=64m
# Reverse-proxied by Caddy/nginx on the host; no public bind here.

[Service]
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
